# systemd/swayr.service
[Unit]
Description=Swayr
After=environment.target

[Service]
ExecStart=/usr/bin/swayrd
KillMode=process
ExecCondition=sh -c '[[  -z ${DISPLAY} && "$XDG_SESSION_TYPE" != "wayland" ]] && exit -1 || exit 0'
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
